# dot/.anacron/cron.daily/backup
#!/bin/bash
######
~/bin/loc_back -l -b 
back1=~/.priv
back2=$back1/back2
rdiff-backup  ~/.config/priv1 $back2
rdiff-backup --remove-older-than 1W $back2
run-parts $back1

