# dot/.anacron/cron.daily/backup
#!/bin/bash
######
cd ~
~/bin/loc_back.sh -l -v 
priv1=~/.priv
run-parts $priv1

