# dot/.anacron/cron.weekly/backup
#!/bin/bash
priv1=~/.priv
priv2=/media/backup/`whoami`/

ls $priv2
sleep 3
rsync -avz --delete $priv1 $priv2
