# dot/.anacron/cron.weekly/signimg
#!/bin/bash
~/bin/loc_back.sh  -b  
priv1=~/.priv
backup=/media/backup/`whoami`
priv2=$backup/

ls $priv2
sleep 3
rsync -avz --delete $priv1 $priv2
